# snippets/supervisor.conf
# http://supervisord.org
[program:queue]
command=/usr/bin/example --loglevel=%(ENV_LOGLEVEL)s
process_name=%
